# Base-OS-Home Automation/Projet/SystemD/HTS-service.service
[Unit]
Description=HTS- Temperature et Humidite - Agent

[Service]
Type=idle
ExecStart=/usr/libexec/HTS

[Install]
WantedBy=multi-user.target

